# provisioning/roles/common/templates/ansible-update.service
[Unit]
Description="Ansible update"

[Service]
ExecStart=ansible-pull
